# templates/profile.d/airflow.j2
export AIRFLOW_HOME={{ airflow_config.core.home }}
